# syncthing/.config/systemd/user/syncthing.service.d/override.conf
[Service]
ExecStart=
ExecStart=/run/current-system/sw/bin/syncthing
